# configs/supervisor/newmanga_server.conf
[program:newmanga_server]
command=/home/akiraak/projects/newmanga/server/run_server_pro.sh
user=akiraak
autorestart=true
stderr_logfile=/var/log/newmanga_server.log